# files/systemd/system/auto-update-flatpak.service
[Unit]
Description = Auto update flatpak

[Service]
ExecStart = /usr/bin/flatpak upgrade -y
