# files/systemd/system/cockpit.socket.d/override.conf
[Socket]
ListenStream=
ListenStream=9191
